# rootfs/etc/profile.d/envvars.sh
#!/bin/sh

export PATH=$PATH:/home/scripts:/home/bash-on-windows/dgx80/lx-tools
alias php-ini='/home/bash-on-windows/home/dgx80/lx-tools/php-ini.sh'
